# src/rootfs/etc/services.d/syslogd/run
#!/bin/sh
set -e

exec s6-setuidgid root syslogd -O /dev/stderr -S -n
