# src/rootfs/etc/services.d/syslog-stdout/run
#!/usr/bin/with-contenv /bin/sh

/usr/sbin/syslog-stdout
