# src/rootfs/etc/services.d/syslog-stdout/run
#!/bin/sh
set -e

exec s6-setuidgid root syslog-stdout
