# src/rootfs/etc/services.d/nginx/run
#!/bin/sh
set -e

s6-svwait /var/run/s6/services/syslog-stdout

exec s6-setuidgid root nginx -g 'daemon off;'
