# src/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv /bin/sh

/usr/sbin/nginx
